# .config/crontab/config
*/1 * * * * /home/dkdk/sh/git/git_config.sh
