# .config/crontab/config
*/10 * * * * /home/dkdk/sh/git/git_config.sh
*/10 * * * * /home/dkdk/sh/git/git_code.sh
*/10 * * * * /home/dkdk/sh/git/git_kernel.sh
** 1 * * * su - root rdate -s time.bora.net
